# root/etc/cron.d/nethserver-service-monitor
# Monitor services status every 5 minutes
*/5 * * * * /usr/sbin/nethserver-service-monitor
